# sys-apps/mpss-daemon/files/mpss
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# Source MPSS function library.
. /etc/mpss/mpss.common

load()
{
	ebegin "Loading MIC module"

	mic_load_modules

	eend $?
}

start() {
	# Ensure that the driver is loaded
	load

	mic_check
	retval=$?; [ $retval -eq 0 ] || return $retval
	
	ebegin "Starting Intel(R) MPSS"

	start-stop-daemon --start --background --pidfile /var/run/mpss --make-pidfile --name mpssd --exec /bin/bash -- -c "exec $exec -l >> /var/log/mpssd 2>&1"
	RETVAL=$? 
	if [ "$RETVAL" = "0" ]; then
		# wait for possible transition from resetting -> ready
		micctrl -w 1> /dev/null
		# wait for the boot process to start
		sleep 5
		# wait for transition form booting -> online
		micctrl -w 1> /dev/null
		RETVAL=$?
	fi

	eend $RETVAL
}

stop() {
	ebegin "Stopping Intel(R) MPSS"
	start-stop-daemon --stop --signal QUIT --pidfile /var/run/mpss
	RETVAL=$?

	if [ "$RETVAL" = "0" ]; then
		# wait for the shutdown process to start
		sleep 1
		micctrl -w 1> /dev/null
		RETVAL=$?
		if [ "$RETVAL" = "9" ]; then
			echo -n $"Shutting down Intel(R) MPSS by force: "
			micctrl -r 1> /dev/null
			RETVAL=$?
			if [ "$RETVAL" = "0" ]; then
				micctrl -w 1> /dev/null
				RETVAL=$?
			fi
		fi
	fi
	
	eend $RETVAL
}
